# supervisor.conf
[program:videostore]
user=videostore
directory=/home/videostore/videostore
command=/home/videostore/videostore/.venv/bin/gunicorn --config gunicorn.py videostore.wsgi:app
